# supplementary/funchotorn-updater.timer
[Unit]
Description=Runs the funchotorn MLS database updater

[Timer]
OnCalendar=*-*-* 1:00:00 UTC
